# .config/autostart/minergate.desktop
[Desktop Entry]
Type=Application
Hidden=false
Exec=/opt/minergate/minergate --auto
Name=Minergate
Terminal=false
